# django_nginx/templates/django_nginx/sub.domain.ext.socket
[Unit]
Description=gunicorn daemon for {{FQDN}}

[Socket]
ListenStream={{socket_path}}


[Install]
WantedBy=sockets.target